# build-context/scheduler/crontabs/charge-lnd
*/15 * * * * /usr/scheduler/jobs/charge-lnd.sh